# buildtools/debian-rcp100.conf
# rcp100 - IP router for Linux platforms
#
# RCP100 is a modular router control plane for Linux platforms. It can be used
# as a standalone router, or it can be deployed as an embedded router in a
# virtualized environment. 

description	"rcp100 router"
author		"RCP Team <rcpteam@yahoo.com>"

start on runlevel [2345]
stop on runlevel [^2345]

console output

pre-start	exec /opt/rcp/bin/start.sh
post-stop exec /opt/rcp/bin/stop.sh
